# supervisord_sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D             ; the program (relative uses PATH, can take args)
autostart=true                ; start at supervisord start (default: true)
startretries=3                ; max # of serial start failures when starting (default 3)
autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
